# riak/db/riak.conf
[program:riak]
command=/bin/bash /riak/run-riak.sh
autorestart=true
startretries=10