# yugabyte/db/yuga.conf
[program:yb-master]
command=/bin/bash /yuga/run-yb-master.sh
autorestart=true
startretries=10

[program:yb-tserver]
command=/bin/bash /yuga/run-yb-tserver.sh
autorestart=true
startretries=10