# cockroachdb/db/cockroachdb.conf
[program:cockroachdb]
command=/bin/bash /roach/run-roach.sh
autorestart=true
startretries=10