# mongodb/db/mongo.conf
[program:mongodb]
command=/bin/bash /mongo/run-mongo.sh
autorestart=true
startretries=10

[program:api]
command=/bin/bash /mongo/run-api.sh
autorestart=true
startretries=10