# rethinkdb/db/rethink.conf
[program:rethink]
command=/rethink/run-rethink.sh
autorestart=true
startretries=10