# etcd/db/etcd.conf
[program:etcd]
command=/bin/bash /etcd/run-etcd.sh
autorestart=true
startretries=10