# consul/db/consul.conf
[program:consul]
command=/bin/bash /consul/run-consul.sh
autorestart=true
startretries=10